# sh/profile.d/date.lazy.sh
# .-------------------.
# | Show current date |
# '-------------------'

echo "Today is $(LC_ALL=C date)"
