# sh/profile.d/misc.sh
# .---------------.
# | Misc settings |
# '---------------'

ulimit -S -c 0	# core dumpsize
umask 022
mesg n 2>/dev/null

return 0
